# install/systemd/dns-monitor.service
[Unit]
Description=Serviço de Monitoramento DNS para ISPs
After=network.target unbound.service
Wants=unbound.service

[Service]
Type=simple
ExecStart=/opt/dns-protection/scripts/monitoring/dns-monitor.sh
Restart=always
RestartSec=10
User=root
Group=root

[Install]
WantedBy=multi-user.target